# docker/web/common/supervisord-nginx.conf
[program:nginx]
command=service nginx start
priority=900
username=www-data
autorestart=true
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log